# infra/cloudflared/cloudflared-aitoolgeek.service
# Systemd unit for the aitoolgeek tunnel.
# Matches the pattern of cloudflared-pratibedan.service.
#
# DESTINATION on server: /etc/systemd/system/cloudflared-aitoolgeek.service
#
# After placing, run:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cloudflared-aitoolgeek.service
#   sudo systemctl status cloudflared-aitoolgeek.service

[Unit]
Description=Cloudflare Tunnel for AI Tool Geek
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR-USER
ExecStart=/usr/bin/cloudflared tunnel --config /home/YOUR-USER/.cloudflared/config-aitoolgeek.yml run aitoolgeek
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
